# 8.0/etc/supervisor/conf.d/20_odoo.conf
[program:odoo]
command=/usr/bin/openerp-server
user=odoo

