# 8.0/etc/supervisor/conf.d/20_odoo.conf
[program:odoo]
environment=HOME=/var/lib/odoo
command=/usr/bin/openerp-server --config /etc/odoo/openerp-server.conf --logfile /var/log/odoo/odoo-server.log
user=odoo

